# DevOps Engineer/YAT206_GIZ1_SWD1_G2e_DEPI_30Apr/Group 4 (Ultimate CI-CD Pipeline)/DevOpsPipelineMastery-main/ansible/templates/blackbox.service.j2
[Unit]
Description=Black Box Service
After=network.target

[Service]
ExecStart=/opt/prometheus/blackbox/blackbox_exporter --config.file=/opt/prometheus/blackbox//blackbox.yml
Restart=always
User=prometheus 
Group=prometheus

[Install]
WantedBy=multi-user.target
