# DevOps Engineer/YAT206_GIZ1_SWD1_G2e_DEPI_30Apr/Group 4 (Ultimate CI-CD Pipeline)/DevOpsPipelineMastery-main/ansible/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter Service
After=network.target

[Service]
ExecStart=/opt/prometheus/node_exporter/node_exporter
Restart=always
User=prometheus 
Group=prometheus

[Install]
WantedBy=multi-user.target
